# deploy/supervisor/poly15minbtc.conf
[group:poly15minbtc]
programs=poly15minbtc-collector,poly15minbtc-dashboard
priority=200

[program:poly15minbtc-collector]
directory=/home/ubuntu/codexProxy/codex2api/deployments/Poly-15minBTC
command=/home/ubuntu/codexProxy/codex2api/deployments/Poly-15minBTC/scripts/run_collector.sh
autostart=true
autorestart=true
startsecs=5
startretries=10
stopsignal=TERM
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/home/ubuntu/codexProxy/codex2api/deployments/Poly-15minBTC/logs/collector.stdout.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/home/ubuntu/codexProxy/codex2api/deployments/Poly-15minBTC/logs/collector.stderr.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5

[program:poly15minbtc-dashboard]
directory=/home/ubuntu/codexProxy/codex2api/deployments/Poly-15minBTC
command=/home/ubuntu/codexProxy/codex2api/deployments/Poly-15minBTC/scripts/run_dashboard.sh
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1",POLY15_DASHBOARD_PORT="18080"
stdout_logfile=/home/ubuntu/codexProxy/codex2api/deployments/Poly-15minBTC/logs/dashboard.stdout.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/home/ubuntu/codexProxy/codex2api/deployments/Poly-15minBTC/logs/dashboard.stderr.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
